# automated_IUTs/coap_client_coapthon/supervisor.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:agent]
directory= ./coap_testing_tool/agent/
command = sh -c "sleep 1;python -m agent connect --url %(ENV_AMQP_URL)s --name coap_client_agent"
user=root
stopsignal=INT
stopasgroup=true
autorestart=false
loglevel=debug
; for sending logs to docker
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
;redirect_stderr=true
;stdout_logfile = /var/log/agent-coap_client_agent-stdout.log
;stdout_logfile_maxbytes = 10MB
;stdout_logfile_backups = 5

[program:bootstrap-agent]
; auto-bootstrap the agent before automated-iut boots up
command = sh -c "sleep 4;/usr/bin/python3 -m coap_testing_tool.agent.utils.bootstrap_agent coap_client_agent bbbb :1"
autorestart=false
stopsignal=INT
stopasgroup=true
loglevel=debug
; for sending logs to docker
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
;redirect_stderr=true
;stdout_logfile = /var/log/bootstrap_agent-stdout.log

[program:automated-iut]
command = sh -c "sleep 10;/usr/bin/python3 -m automated_IUTs.coap_client_coapthon.automated_iut"
autorestart=false
stopsignal=INT
stopasgroup=true
loglevel=debug
; for sending logs to docker
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
;redirect_stderr=true
;stdout_logfile = /var/log/automated_iut-coap_client_coapthon-stdout.log

[program:automated-iut-ipv6-routes-config]
command = sh -c "sleep 10;ip -6 route add cccc::/64 dev tun0"
autorestart=false
stopsignal=INT
loglevel=debug
; for sending logs to docker
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

